# .config/timers/pacman-index-update/pacman-index-update.service
[Unit]
Description=Update the Pacman package index (for the statusbar)

[Service]
Type=oneshot
ExecStart=pacman -Syq
